# deploy/kraymini.supervisor.conf
; kraymini 的 Supervisor 配置示例
; 用法：复制到 /etc/supervisor/conf.d/kraymini.conf（路径因发行版而异），
;       按需修改 command / user / 日志路径，然后：
;       sudo supervisorctl reread && sudo supervisorctl update
;
; 依赖：已安装 kraymini（pip/venv）且 PATH 中能启动 xray，或 kraymini 配置里指定 xray 绝对路径。

[program:kraymini]
; 显式指定配置文件；未使用 -c 时 kraymini 会按默认路径搜索 config.toml
command=/usr/local/bin/kraymini run -c /etc/kraymini/config.toml
; 若使用虚拟环境，改为例如：command=/opt/kraymini/.venv/bin/kraymini run -c /etc/kraymini/config.toml

directory=/var/lib/kraymini
user=nobody
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=10
stopasgroup=true
killasgroup=true
redirect_stderr=false
stdout_logfile=/var/log/kraymini/supervisor.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
stderr_logfile=/var/log/kraymini/supervisor.err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3

; 若 xray 不在默认 PATH 中，可在此注入（按实际安装路径修改）
environment=PATH="/usr/local/bin:/usr/bin:/bin"
